# crontab
* * * * * $HOME/obd/update
* * * * * $HOME/obd/post-to-hass
