# crontab
PATH=$HOME/obd/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * $HOME/obd/update
* * * * * $HOME/obd/post-to-hass
* * * * * $HOME/obd/battery-check
