# crontab
* * * * * $HOME/obd/ci
* * * * * $HOME/obd/post-to-hass
